# plugins/notification-tools/skills/notification-center/launchd/com.saber.notification-center.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.saber.notification-center.watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/saberrao/.codex/skills/notification-center/watcher.py</string>
  </array>

  <key>StartInterval</key>
  <integer>14400</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/saberrao/vault/.notification-center/.watch.launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/saberrao/vault/.notification-center/.watch.launchd.stderr.log</string>
</dict>
</plist>
